# deploy/lightai-agent.service
[Unit]
Description=LightAI Agent (GPU node)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lightai
Group=lightai
WorkingDirectory=/opt/lightai
ExecStart=/opt/lightai/bin/lightai-agent --config /opt/lightai/lightai-agent.toml
Restart=always
RestartSec=5

# CRITICAL: KillMode=process ensures systemd only terminates the Agent itself,
# not the entire cgroup. Model instance processes started by the Agent must
# survive Agent restarts independently.
# Default KillMode=control-group would kill all child processes on Agent exit.
KillMode=process

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/lightai/data /opt/lightai/logs /opt/lightai/models
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
