# deploy/lightai-server.service
[Unit]
Description=LightAI Server (control plane)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lightai
Group=lightai
WorkingDirectory=/opt/lightai
# Put LIGHTAI_SERVER_CONFIG and optional non-secret runtime settings here.
# First admin creation should use the Web setup flow.
EnvironmentFile=/etc/lightai/lightai-server.env
ExecStart=/opt/lightai/bin/lightai-server --config /opt/lightai/lightai-server.toml
Restart=always
RestartSec=3

# Server manages no child processes; KillMode=process or control-group are both safe.
# Explicitly set to process for consistency with Agent config.
KillMode=process

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/lightai/data /opt/lightai/logs
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
